# etc/systemd/system/evremap-restart.service
[Unit]
Description=Restart evremap when a keyboard is added

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart evremap.service
